# services/neural_infer.service
[Unit]
Description=Neural Inference Service (runs inside ros2 container)
PartOf=sim-session.service
After=sim-session.service neural_executor.service
Wants=neural.target
Before=neural.target

[Service]
Type=simple
# Dynamic container discovery - find the running ros2 container
ExecStartPre=/bin/sh -c 'CONTAINER=$$(docker ps --filter "name=ros2" --format "{{.Names}}" | head -n 1); \
if [ -z "$$CONTAINER" ]; then \
  echo "Error: ros2 container not found" >&2; \
  exit 1; \
fi; \
echo "$$CONTAINER"'
# Execute neural_infer.py inside the ros2 container
# Using container-local Python intentionally
ExecStart=/bin/sh -c 'CONTAINER=$$(docker ps --filter "name=ros2" --format "{{.Names}}" | head -n 1); \
docker exec -u ros "$$CONTAINER" /bin/bash -lc "source /opt/ros/humble/setup.bash && \
cd /home/ros/ros2_ws && \
source install/setup.bash && \
python3 src/neural_manager/neural_inference/neural_infer.py"'
Restart=on-failure
RestartSec=5

# Security settings for user service
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
