# services/neural_executor.service
[Unit]
Description=Neural Executor Service (runs inside ros2 container)
After=network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
# Dynamic container discovery - find the running ros2 container
ExecStartPre=/bin/sh -c 'CONTAINER=$$(docker ps --filter "name=ros2" --format "{{.Names}}" | head -n 1); \
if [ -z "$$CONTAINER" ]; then \
  echo "Error: ros2 container not found" >&2; \
  exit 1; \
fi; \
echo "$$CONTAINER"'
# Execute neural_executor launch file inside the ros2 container
ExecStart=/bin/sh -c 'CONTAINER=$$(docker ps --filter "name=ros2" --format "{{.Names}}" | head -n 1); \
docker exec -u ros "$$CONTAINER" /bin/bash -lc "source /opt/ros/humble/setup.bash && \
cd /home/ros/ros2_ws && \
source install/setup.bash && \
ros2 launch neural_executor neural_executor.launch.py"'
Restart=on-failure
RestartSec=5

# Security settings for user service
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
