# services/neural_infer.service
[Unit]
Description=Neural Inference Service (runs inside bht container)

BindsTo=bht-liveness.service
After=bht-liveness.service neural_executor.service
PartOf=sim-session.service
Before=neural.target

[Service]
Type=simple
# Dynamic container discovery - find the running bht container
ExecStartPre=/bin/sh -c 'CONTAINER=$$(docker ps --filter "name=bht" --format "{{.Names}}" | head -n 1); \
if [ -z "$$CONTAINER" ]; then \
  echo "Error: bht container not found" >&2; \
  exit 1; \
fi; \
echo "$$CONTAINER"'
# Execute neural_infer.py inside the bht container
# Using container-local Python intentionally
ExecStart=/bin/sh -c 'CONTAINER=$$(docker ps --filter "name=bht" --format "{{.Names}}" | head -n 1); \
docker exec -u ros "$$CONTAINER" /bin/bash -lc "source /opt/ros/humble/setup.bash && \
cd /home/ros/ros2_ws && \
source install/setup.bash && \
python3 src/neural_manager/neural_inference/neural_infer.py"'
Restart=no

# Security settings for user service
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=neural.target
