# services/sim-session.service
[Unit]
Description=VTOL Simulation Session (tmux + docker compose)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=forking
WorkingDirectory=/home/rec/server/vtol-interface

# Start tmux session with docker compose
ExecStart=/usr/bin/tmux new-session -d -s vtol-sim -x 200 -y 50 \
  'docker compose up px4 qgc && docker compose attach ros2; docker compose down; tmux kill-session -t vtol-sim'

# Wait for containers to be ready, then start default group (neural)
ExecStartPost=/bin/bash -c 'sleep 8 && systemctl --user start neural.target'

# Stop: kill tmux session (which triggers docker compose down)
ExecStop=/usr/bin/tmux kill-session -t vtol-sim 2>/dev/null || true
ExecStopPost=/bin/bash -c 'docker compose down 2>/dev/null || true'

# When sim-session stops, stop all dependent services
ExecStopPost=/bin/bash -c 'systemctl --user stop neural.target test.target 2>/dev/null || true'

Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
