# services/test.target
[Unit]
Description=Test Tools Group B
Conflicts=neural.target
AllowIsolate=yes

[Install]
WantedBy=sim-session.service
