# services/test.target
[Unit]
Description=Test Executor Group B (test_executor with joystick)
Conflicts=neural.target

[Install]
WantedBy=sim-session.service
